# scripts/daily-news-web.service
[Unit]
Description=Odoo Daily News Web UI + Scheduler
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/daily-news
EnvironmentFile=/path/to/daily-news/.env
ExecStart=/path/to/daily-news/.venv/bin/python app.py
Restart=on-failure
RestartSec=5
User=www-data

[Install]
WantedBy=multi-user.target
